# cron/crontab
# Homelab crontab — managed by cron/install.sh
# m h  dom mon dow   command

*/5 * * * * /home/micjcameron/homelab/home-automation/scripts/lifeline.sh >> /home/micjcameron/lifeline.log 2>&1
0 10 * * *  /home/micjcameron/homelab/home-automation/scripts/lifeline_daily.sh >> /home/micjcameron/lifeline.log 2>&1
0 7,13,17 * * * /home/micjcameron/homelab/personal-scripts/parent-reminder/parent-reminder.sh >> /home/micjcameron/parent-reminder.log 2>&1
0 20 * * *      /home/micjcameron/homelab/personal-scripts/bin-reminder/bin-reminder.sh tomorrow >> /home/micjcameron/bin-reminder.log 2>&1
0 8 * * *       /home/micjcameron/homelab/personal-scripts/warmup-reminder/warmup-reminder.sh >> /home/micjcameron/warmup-reminder.log 2>&1
45 8 * * *      /home/micjcameron/homelab/personal-scripts/bin-reminder/bin-reminder.sh today    >> /home/micjcameron/bin-reminder.log 2>&1
0 4 1 * *       /home/micjcameron/homelab/home-automation/scripts/docker-prune.sh >> /home/micjcameron/docker-prune.log 2>&1
